# files/etc/hotplug.d/iface/95-wan3-manager
#!/bin/sh
# Re-evaluate local WANs and tunnel policy after netifd events.

[ "$(uci -c /etc/config.mesh -q get aredn.multiwan.enabled 2>/dev/null)" = 1 ] || exit 0

case "${INTERFACE:-}" in
    wan|wan2|wan3)
        /usr/local/bin/wan-route-cache "$INTERFACE" >/dev/null 2>&1 || true
        /usr/local/bin/wan3-manager route-event "$INTERFACE" "${ACTION:-}" >/dev/null 2>&1 &
        ;;
    wg*|tun*)
        /usr/local/bin/wan-tunnel-guard apply >/dev/null 2>&1 &
        ;;
esac
